# arkhen-os/config/mcp-server.service
[Unit]
Description=Arkhe MCP Server Socket

[Socket]
ListenStream=/run/mcp-server/mcp.sock
User=arkhe
Group=arkhe

[Install]
WantedBy=sockets.target
